# user/systemd/kb-layout-watcher.service
[Unit]
Description=Keyboard layout auto-switcher (GUI)
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/kb-layout-watcher.sh
Restart=on-failure

[Install]
WantedBy=default.target
